# roles/client/templates/web-client.service.j2
[Unit]
Description=Caliopen web client

[Service]
WorkingDirectory=/opt/caliopen/client
ExecStart=/opt/caliopen/client/bin/server --config config/server.json

[Install]
WantedBy=multi-user.target
