# roles/client/templates/web-client.service.j2
[Unit]
Description=Caliopen web client

[Service]
Environment
WorkingDirectory=/opt/caliopen/server
ExecStart=/usr/bin/node index.js
[Install]
WantedBy=multi-user.target