# roles/client/templates/web-client.service.j2
[Unit]
Description=Caliopen web client

[Service]
Environment
WorkingDirectory=/opt/caliopen/client
ExecStart=start-client

[Install]
WantedBy=multi-user.target