# roles/object-store/templates/minio.service.j2
[Unit]
Description=Minio object store service

[Service]
Restart=always
ExecStart=/usr/local/bin/minio server --config-dir /etc/minio /srv/data
ExecStop=pkill minio

[Install]
WantedBy=local.target
