# roles/smtp/files/etc/cron.d/spamassassin
# Spamassassin care
42 */12 * * * root /root/sa-care.sh
